# miraikakakubatch/config/miraikakaku-batch.service
[Unit]
Description=Miraikakaku Batch System
After=network.target
Wants=network.target

[Service]
Type=simple
User=yuuku
WorkingDirectory=/mnt/c/Users/yuuku/cursor/miraikakaku
ExecStart=/usr/bin/python3 /mnt/c/Users/yuuku/cursor/miraikakaku/miraikakakubatch.py --mode monitor
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
